# deploy/systemd/robominer-web.service
[Unit]
Description=RoboMiner Rust web host
Documentation=file:/opt/robominer/deploy/systemd/README.md
After=network-online.target mariadb.service mysql.service
Wants=network-online.target

[Service]
Type=simple
User=robominer
Group=robominer
WorkingDirectory=/opt/robominer
ExecStart=/opt/robominer/bin/robominer-web --config /etc/robominer/robominer.conf
# Fail the unit if loopback /health does not become ready (DB + migrations).
ExecStartPost=/opt/robominer/bin/robominer-wait-web-health /etc/robominer/robominer.conf
Restart=on-failure
RestartSec=10
TimeoutStartSec=90
TimeoutStopSec=30
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
